# automations/launchd/com.carlos.brain.weekly-synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.carlos.brain.weekly-synthesis</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>/Users/carlos/Brain/automations/scripts/run-weekly-synthesis.sh</string>
  </array>

  <!-- Fires Monday at 09:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StartCalendarIntervalRunOnWake</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/carlos/Brain</string>

  <key>StandardOutPath</key>
  <string>/Users/carlos/Brain/automations/scripts/logs/launchd-weekly-synthesis.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/carlos/Brain/automations/scripts/logs/launchd-weekly-synthesis.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
